# spark/container/root/etc/services.d/java/run
#!/bin/bash -e

source /scripts/with-bigcontenv

# Start java app
exec s6-setuidgid pi java -version
